# springboot-gradle/dist/spring-boot-demo.service
## Instruction
# 1. Install the spring boot application in /var/springboot-gradle-demo, 
# 2. Copy this script named spring-boot-demo.service in /etc/systemd/system directory
[Unit]
Description=springboot-gradle-demo
After=syslog.target

[Service]
User=springBootDemo
ExecStart=/var/springboot-gradle-demo/springboot-gradle-demo.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target